# sample-service/sample-service.service
[Unit]
Description=Sample HTTP Service
Documentation=https://github.com/your-repo/sample-service
After=network.target

[Service]
Type=simple
User=service-manager
Group=service-manager
WorkingDirectory=/opt/service-manager/sample-service
ExecStart=/opt/service-manager/sample-service/sample-service.py
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sample-service

# Environment variables
Environment=PYTHONUNBUFFERED=1

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/service-manager/sample-service

[Install]
WantedBy=multi-user.target
